# roles/platzigram/files/upstart/platzigram-api-pictures.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid root
chdir /opt/platzigram-api

env NODE_ENV="production"
env PLATZIGRAM_SECRET="5ec7b912-6634-4db1-9fbf-1f4717fbed1e"

exec npm run pictures
